# .deploy/cron/crontab
@reboot sleep 300 && cd /var/www/html && bin/console cache:clear
*/2 * * * * echo "cron works" > /tmp/cron.ping
